# mac/launchd/local.rgnrd.audio-analyzer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.rgnrd.audio-analyzer</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>export PATH="/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:$HOME/.local/bin"; cd "$HOME/code/reginerd-fm" &amp;&amp; source .env &amp;&amp; mkdir -p "$HOME/Library/Logs/rgnrd" &amp;&amp; exec uv run python mac/audio_analyzer.py &gt;&gt; "$HOME/Library/Logs/rgnrd/audio-analyzer.log" 2&gt;&amp;1</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Users/reggiecooper/.local/bin</string>
    <key>HOME</key>
    <string>/Users/reggiecooper</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/reggiecooper/Library/Logs/rgnrd/audio-analyzer.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/reggiecooper/Library/Logs/rgnrd/audio-analyzer.log</string>
</dict>
</plist>
